# docker/backups/crontab
#m	 h	  d	   m    w	 command
0    1    *    *    *    sh /etc/periodic/backups/postgresql
